# ai-platform/browser-container/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

# Lets `docker exec ai-browser supervisorctl status` report per-service health.
[unix_http_server]
file=/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

# Virtual framebuffer — everything else draws into display :0.
[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1920x1080x24
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

# XFCE desktop for user claude. dbus-launch gives it the session bus it needs to
# render the panel/window manager (without it noVNC shows a black screen). Waits
# for Xvfb to accept connections first so the session doesn't start against a
# missing display and get thrown into a restart loop.
[program:xfce]
command=/bin/bash -c 'until xdpyinfo -display :0 >/dev/null 2>&1; do sleep 0.5; done; exec dbus-launch --exit-with-session startxfce4'
user=claude
environment=HOME="/home/claude",USER="claude",DISPLAY=":0"
autorestart=true
startsecs=3
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

# VNC server. Must not start until X is actually up, otherwise x11vnc exits with
# "can't open display" and the port never opens — which is the "Failed to connect
# to server" noVNC error.
[program:x11vnc]
command=/bin/bash -c 'until xdpyinfo -display :0 >/dev/null 2>&1; do sleep 0.5; done; exec x11vnc -display :0 -forever -shared -rfbport 5900 -passwd "%(ENV_VNC_PASSWORD)s" -noxdamage'
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

# noVNC/websockify: serves vnc.html on 6080 and bridges the browser websocket to
# x11vnc on 5900. Reachable from the host as :9901 (published by proxy-gateway).
[program:novnc]
command=/opt/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
